# data-collection/crontab
* * * * * python3 read_sensors.py >> data/sensors.csv